# support/cron.d
0 18 * * * horus cd /home/horus/bin && ./grab.pl --quiet --email=dcops@fusionone.com
